# com.philreviews.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.philreviews.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mzwolinski/PhilReview/update.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mzwolinski/PhilReview</string>

    <!-- Run every Sunday at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mzwolinski/PhilReview/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mzwolinski/PhilReview/logs/launchd_stderr.log</string>

    <!-- Don't run multiple instances -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
